# blazegraph-rpm/src/rpm/scripts/preremove
service blazegraph stop
chkconfig blazegraph off
